# unban-bot.service
[Unit]
Description=Unban Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/unban
EnvironmentFile=/root/unban/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /root/unban/index.js
Restart=always
RestartSec=10
StandardOutput=append:/root/unban/logs/bot.log
StandardError=append:/root/unban/logs/error.log

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
